# ipfs-accelerate-mcp.service
[Unit]
Description=IPFS Accelerate MCP Server with Dashboard
Documentation=https://github.com/endomorphosis/ipfs_accelerate_py
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=barberb
Group=barberb
WorkingDirectory=/home/barberb/ipfs_accelerate_py

# Environment variables
Environment="PATH=/home/barberb/ipfs_accelerate_py/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="MCP_HOST=0.0.0.0"
Environment="MCP_PORT=9000"

# Note: GitHub tokens should be configured via `gh auth login` instead of environment variables
# This allows the GitHub CLI to manage authentication securely

# Start the MCP server using the CLI tool
ExecStart=/home/barberb/ipfs_accelerate_py/.venv/bin/ipfs-accelerate mcp start \
    --dashboard \
    --host 0.0.0.0 \
    --port 9000 \
    --keep-running

# Restart configuration
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ipfs-accelerate-mcp

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
MemoryMax=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
