# ipfs-accelerate-mcp.service
[Unit]
Description=IPFS Accelerate MCP Server with Dashboard
Documentation=https://github.com/endomorphosis/ipfs_accelerate_py
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=devel
Group=devel
WorkingDirectory=/home/devel/ipfs_accelerate_py

# Environment variables
Environment="PATH=/home/devel/ipfs_accelerate_py/ipfs_env/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/devel/ipfs_accelerate_py"
Environment="VIRTUAL_ENV=/home/devel/ipfs_accelerate_py/ipfs_env"
Environment="MCP_HOST=0.0.0.0"
Environment="MCP_PORT=9000"

# Update before starting
ExecStartPre=/home/barberb/ipfs_accelerate_py/scripts/auto-update.sh

# Start the MCP server using wrapper script
ExecStart=/home/devel/ipfs_accelerate_py/start_mcp_server.sh

# Restart configuration
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ipfs-accelerate-mcp

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
MemoryMax=4G
CPUQuota=200%
LimitNOFILE=8192

[Install]
WantedBy=multi-user.target
