# deployments/systemd/ipfs-accelerate.service
[Unit]
Description=IPFS Accelerate (primary) MCP + TaskQueue P2P service
Documentation=https://github.com/endomorphosis/ipfs_accelerate_py
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=devel
Group=devel
WorkingDirectory=%h/ipfs_accelerate_py
Environment=VIRTUAL_ENV=%h/ipfs_accelerate_py/.venv
Environment=PATH=%h/ipfs_accelerate_py/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=%h/ipfs_accelerate_py

# Optional secrets/env overrides (recommended):
#   GH_TOKEN=... (preferred for gh)
#   GITHUB_TOKEN=... (optional)
#   IPFS_ACCELERATE_GITHUB_REPO=owner/repo
EnvironmentFile=-/etc/ipfs-accelerate/secrets.env

# If you use `gh auth login` (user-scoped auth), point gh at the service user's
# config directory. This is especially important when the unit is installed with
# `deployments/systemd/install.sh --user <you>`.
Environment=GH_CONFIG_DIR=%h/.config/gh

# Keep caches/state out of $HOME (ProtectHome=read-only)
CacheDirectory=ipfs-accelerate
StateDirectory=ipfs-accelerate
LogsDirectory=ipfs-accelerate
RuntimeDirectory=ipfs-accelerate
Environment=PIP_CACHE_DIR=/var/cache/ipfs-accelerate/pip
Environment=IPFS_ACCELERATE_STATE_DIR=/var/lib/ipfs-accelerate
Environment=HOME=/var/lib/ipfs-accelerate
Environment=XDG_CACHE_HOME=/tmp
Environment=XDG_STATE_HOME=/var/lib/ipfs-accelerate
Environment=XDG_CONFIG_HOME=/var/lib/ipfs-accelerate/config

# P2P Cache Configuration (libp2p)
Environment=CACHE_ENABLE_P2P=false
Environment=CACHE_LISTEN_PORT=9100
Environment=IPFS_ACCELERATE_GITHUB_REPO=endomorphosis/ipfs_accelerate_py

# TaskQueue P2P RPC service + worker autoscaling are enabled by default by the MCP CLI.
# Use env overrides only if you need to disable or tune behavior.

# Batching + multi-peer draining
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_MESH_PEER_FANOUT=1
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_MESH_CLAIM_BATCH=1
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_LOCAL_CLAIM_BATCH=1

# Enable additional HuggingFace task types (beyond text-generation)
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_ENABLE_HF=1
# Minimal HF fallback is opt-in for this integrated service.
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_MINIMAL_HF=0

# Per-task micro-batch caps (kept disabled by default here)
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_TEXTGEN_BATCH_MAX=1
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_TEXT2TEXT_BATCH_MAX=1
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_TEXTCLS_BATCH_MAX=1
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_EMBED_BATCH_MAX=1

# Estimator tuning
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_HF_MEM_FRACTION=0.75
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_HF_MEM_RESERVE_MB=1024
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_HF_ACT_MULT=4.0

Environment=IPFS_ACCELERATE_PY_TASK_WORKER_ENABLE_DOCKER=1
Environment=IPFS_DATASETS_PY_TASK_WORKER_ENABLE_DOCKER=1
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_ENABLE_COPILOT_CLI=1
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_TASK_TYPES=tool.call,text-generation,text2text-generation,embedding,text-classification,hf.pipeline,docker.execute,docker.execute_docker_container,docker.github
Environment=IPFS_DATASETS_PY_TASK_WORKER_TASK_TYPES=tool.call,text-generation,text2text-generation,embedding,text-classification,hf.pipeline,docker.execute,docker.execute_docker_container,docker.github
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_ALLOW_COPILOT_CONTINUE_WITHOUT_RESUME=1
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_MESH_REFRESH_S=1.0
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_MESH_CLAIM_INTERVAL_S=0.1
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_MESH_MAX_PEERS=25
Environment=IPFS_ACCELERATE_PY_TASK_P2P_ENABLE_TOOLS=1
Environment=IPFS_ACCELERATE_PY_TASK_P2P_ENABLE_CACHE=1
Environment=IPFS_ACCELERATE_PY_TASK_P2P_LISTEN_PORT=9100
Environment=IPFS_DATASETS_PY_TASK_P2P_LISTEN_PORT=9100
Environment=IPFS_ACCELERATE_PY_TASK_QUEUE_PATH=/var/lib/ipfs-accelerate/task_queue.duckdb
Environment=IPFS_ACCELERATE_PY_TASK_P2P_CACHE_DIR=/var/cache/ipfs-accelerate/p2p_cache
# Auto-detect a non-loopback IP for the announce multiaddr
Environment=IPFS_ACCELERATE_PY_TASK_P2P_PUBLIC_IP=auto
# Make announce file visible outside the PrivateTmp sandbox
Environment=IPFS_ACCELERATE_PY_TASK_P2P_ANNOUNCE_FILE=/var/cache/ipfs-accelerate/task_p2p_announce.json

# Enable /dnsaddr/... -> concrete multiaddr expansion for public bootstraps
Environment=IPFS_ACCELERATE_PY_TASK_P2P_DNSADDR_RESOLVE=1
Environment=IPFS_DATASETS_PY_TASK_P2P_DNSADDR_RESOLVE=1

# Updates + pip installs can take longer than systemd's default 90s
TimeoutStartSec=0

ExecStart=%h/ipfs_accelerate_py/.venv/bin/python3 -m ipfs_accelerate_py.mcp.cli --host 0.0.0.0 --port 9000 --log-level INFO
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ipfs-accelerate

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=%h/ipfs_accelerate_py /var/cache/ipfs-accelerate /var/lib/ipfs-accelerate /var/log/ipfs-accelerate /run/ipfs-accelerate
ProtectHome=read-only

[Install]
WantedBy=multi-user.target