# deployments/systemd/ipfs-accelerate-task-worker.service
[Unit]
Description=ipfs_accelerate_py Task Worker (P2P + mDNS mesh)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Install script patches User/Group/paths when using deployments/systemd/install.sh --user <you>
User=devel
Group=devel

# Adjust to your checkout or installed location.
WorkingDirectory=%h/ipfs_accelerate_py
Environment=VIRTUAL_ENV=%h/ipfs_accelerate_py/.venv
Environment=PATH=%h/ipfs_accelerate_py/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=%h/ipfs_accelerate_py

# Load per-host settings (listen port, session tag, caches, etc.)
EnvironmentFile=-/etc/ipfs-accelerate/task-worker.env

# Enable docker task types by default (override in task-worker.env if needed).
Environment=IPFS_ACCELERATE_PY_TASK_WORKER_ENABLE_DOCKER=1
Environment=IPFS_DATASETS_PY_TASK_WORKER_ENABLE_DOCKER=1

# Session failover (conservative): resubmit locally if a queued task requires a
# different session_id and the matching session peer can't drain it.
# - Enabled by default in code (set IPFS_ACCELERATE_PY_TASK_P2P_SESSION_FAILOVER=0 to disable)
Environment=IPFS_ACCELERATE_PY_TASK_P2P_SESSION_FAILOVER_AFTER_S=240
Environment=IPFS_DATASETS_PY_TASK_P2P_SESSION_FAILOVER_AFTER_S=240

# Run the TaskQueue worker daemon.
# Default behavior includes:
# - hosting the TaskQueue libp2p service
# - mDNS mesh draining
# - autoscaling by spawning/retiring child worker processes
#
# /etc/ipfs-accelerate/task-worker.env remains optional for per-host overrides
# (ports, session tags, queue path, max workers, disabling remote autoscale, etc.).
ExecStart=%h/ipfs_accelerate_py/.venv/bin/python -m ipfs_accelerate_py.p2p_tasks.worker

Restart=always
RestartSec=2

# You almost certainly want a higher limit for model downloads.
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
